# confluent/runit/kafka/finish
#!/bin/sh

sv -w 10 stop /etc/service/schema-registry
